# deploy/sin.service
[Unit]
Description=SIn Authentication Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sin
Group=sin
EnvironmentFile=/etc/sin/config.env
ExecStart=/usr/local/bin/sin-server
Restart=on-failure
RestartSec=5s
WorkingDirectory=/opt/sin

[Install]
WantedBy=multi-user.target
